# apps/posting/deploy/systemd/alexgetman-site-feed.service
[Unit]
Description=AlexGetman site feed worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/opt/alexgetman-posting
EnvironmentFile=/opt/alexgetman-posting/site-feed.env
ExecStart=/usr/bin/docker compose -f /opt/alexgetman-posting/compose.yaml up site-feed
ExecStop=/usr/bin/docker compose -f /opt/alexgetman-posting/compose.yaml stop site-feed
Restart=on-failure
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
