# deploy/alexgetman-deploy-agent.service
[Unit]
Description=AlexGetman immutable image deploy agent
After=docker.service network-online.target
Requires=docker.service

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/home/deploy/repos/alexgetman.com
EnvironmentFile=/etc/alexgetman/deploy-agent.env
ExecStart=/usr/local/bin/bun /home/deploy/repos/alexgetman.com/deploy/deploy-agent.ts
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/lib/alexgetman-deploy /home/deploy/alexgetman-runtime /home/deploy/maru /var/run/docker.sock

[Install]
WantedBy=multi-user.target
