# deploy/media-processor/studio-media-tunnel.service
[Unit]
Description=Studio Media Processor reverse tunnel to PS529
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=alex
ExecStart=/usr/bin/ssh -N -i /home/alex/.ssh/studio-media-processor_ed25519 -o BatchMode=yes -o ExitOnForwardFailure=yes -o ServerAliveInterval=30 -o ServerAliveCountMax=3 -o StrictHostKeyChecking=accept-new -R 127.0.0.1:9087:127.0.0.1:9087 deploy@5.129.238.194
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
